# .config/systemd/user/update-pass-git.service
[Unit]
Description=Do a pass git sync

[Service]
Type=oneshot
ExecStart=pass git pull --rebase=true
ExecStart=pass git push

