# .config/systemd/user/update-dotfiles-git.timer
[Unit]
Description=Run on boot and each 30 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
