# .config/systemd/user/update-pass-git.service
[Unit]
Description=Do a pass git pull

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=pass git pull

